# firmware/Robotino3/tags/1.2.3/firmware/linux_install/xenial/etc/systemd/system/robotino3-firmware.service
[Unit]
Description="robotino3-firmware"

[Service]
ExecStart=/usr/local/robotino/firmware/install.sh
Restart=no

[Install]
