# daemons/scripts/Robotino3/xenial/etc/systemd/system/controld3.service
[Unit]
Description="controld3"

[Service]
ExecStart=/usr/local/robotino/daemons/bin/controld3
Restart=on-failure

[Install]
