# daemons/scripts/Robotino3/xenial/etc/systemd/system/simd.service
[Unit]
Description="simd"

[Service]
ExecStart=/usr/local/robotino/daemons/bin/simd -v
Restart=no

